# deploy/systemd/watchtower-api.service
[Unit]
Description=Watchtower API (Flask backend via Gunicorn)
After=network.target

[Service]
Type=simple
# EDIT these two paths to match your actual deployment location and user
User=watchtower
WorkingDirectory=/opt/recon-ecosystem/watchtower

# Loads MONGO_URI, WATCHTOWER_API_TOKEN, TELEGRAM_*, etc.
EnvironmentFile=/opt/recon-ecosystem/watchtower/.env

ExecStart=/opt/recon-ecosystem/watchtower/venv/bin/gunicorn \
    -w 4 \
    -b 127.0.0.1:3131 \
    --timeout 120 \
    app:app

Restart=on-failure
RestartSec=5

# Basic hardening (relax if it breaks something you need)
NoNewPrivileges=true
PrivateTmp=true

StandardOutput=append:/var/log/watchtower/api.log
StandardError=append:/var/log/watchtower/api.error.log

[Install]
WantedBy=multi-user.target