# deploy/systemd/watchtower-api.service
[Unit]
Description=Watchtower API (Flask backend - Local development/always-on)
After=network.target

# ==============================================================================
# مهم: این سرویس فقط مخصوص نصب روی ماشین شخصی/لوکال (همیشه روشن) است.
# روی VPS ساعتی/موقت نیازی به راه‌اندازی و اجرای این سرویس نیست.
# ==============================================================================

[Service]
Type=simple
# نام کاربری لوکال خود را جایگزین <user> کنید
User=arian
WorkingDirectory=/home/arian/Recon_ecosystem/watchtower

# بارگذاری متغیرهای محیطی لوکال (مثل MONGO_URI و WATCHTOWER_API_TOKEN)
EnvironmentFile=/home/arian/Recon_ecosystem/watchtower/.env

# اجرای مستقیم برای مشاهده لاگ‌ها در حالت Dev/عادی بدون پیچیدگی‌های Gunicorn
ExecStart=/home/arian/Recon_ecosystem/watchtower/venv/bin/python app.py

Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

# ارسال لاگ‌ها به journald (برای بررسی ساده‌تر با journalctl)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target