# deploy/systemd/watchtower-frontend.service
[Unit]
Description=Watchtower Frontend (Node server - Local development/always-on)
After=network.target watchtower-api.service

# ==============================================================================
# مهم: این سرویس فقط مخصوص نصب روی ماشین شخصی/لوکال (همیشه روشن) است.
# روی VPS ساعتی/موقت نیازی به راه‌اندازی و اجرای این سرویس نیست.
# ==============================================================================

[Service]
Type=simple
# نام کاربری لوکال خود را جایگزین <user> کنید
User=arian
WorkingDirectory=/home/arian/Recon_ecosystem/watchtower-frontend

# بارگذاری متغیرهای محیطی فرانت‌اند
EnvironmentFile=/home/arian/Recon_ecosystem/watchtower-frontend/.env

Environment=NODE_ENV=production
Environment=PORT=3000

# اجرای سرور Node.js از روی بیلد فرانت‌اند
ExecStart=/usr/bin/node /home/arian/Recon_ecosystem/watchtower-frontend/dist/index.js

Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

# ارسال لاگ‌ها به journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target